# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
user=root

[program:postgres]
# Wait for initdb (run by nextjs-starter/healthcheck.sh on first boot) before starting
command=sh -c "until [ -f /var/lib/postgresql/data/PG_VERSION ]; do sleep 1; done; exec su postgres -c 'postgres -D /var/lib/postgresql/data'"
autorestart=true
priority=10
stdout_logfile=/var/log/postgres.log
stderr_logfile=/var/log/postgres.log

[program:nextjs-starter]
command=/docker/healthcheck.sh
priority=20
autorestart=false
startsecs=0
stdout_logfile=/var/log/nextjs-starter.log
stderr_logfile=/var/log/nextjs-starter.log
